# deploy/campaignhub-tunnel.service
# =============================================================================
# Campaign Hub — Cloudflare Tunnel Service
# =============================================================================
# Runs a Cloudflare "quick tunnel" so players can access Campaign Hub
# from outside your local network. This gives you a random
# *.trycloudflare.com URL each time it starts.
#
# The URL is logged to the system journal AND written to a file so the
# deploy script can retrieve it easily.
#
# Install on the Pi:
#   sudo cp deploy/campaignhub-tunnel.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable campaignhub-tunnel
#   sudo systemctl start campaignhub-tunnel
#
# Check the tunnel URL:
#   cat /home/rgriffin/campaign-hub/.tunnel-url
#   # or from your Mac:
#   ssh pi "cat ~/campaign-hub/.tunnel-url"
#
# To upgrade to a permanent URL (same address every time), see:
#   https://developers.cloudflare.com/cloudflare-one/connections/connect-networks/
# =============================================================================

[Unit]
Description=Campaign Hub Cloudflare Tunnel
# Start the tunnel after the main app is running
After=network-online.target campaignhub.service
Wants=network-online.target

[Service]
Type=simple
User=rgriffin
Group=rgriffin
WorkingDirectory=/home/rgriffin/campaign-hub

# Run cloudflared and pipe its output through a script that captures the URL.
# cloudflared prints the tunnel URL to stderr, so we redirect stderr to stdout
# and use grep/tee to capture it.
ExecStart=/bin/bash -c '/usr/local/bin/cloudflared tunnel --url http://localhost:3001 2>&1 | while IFS= read -r line; do echo "$line"; if echo "$line" | grep -oP "https://[a-z0-9-]+\\.trycloudflare\\.com" > /dev/null 2>&1; then echo "$line" | grep -oP "https://[a-z0-9-]+\\.trycloudflare\\.com" > /home/rgriffin/campaign-hub/.tunnel-url; fi; done'

# Restart if it crashes (tunnels can drop occasionally)
Restart=on-failure
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=campaignhub-tunnel

[Install]
WantedBy=multi-user.target
